# templates/rc.local
#!/bin/bash

# This is disrupting lscpu on ubuntu 18
# disable hyperthreading
# echo off > /sys/devices/system/cpu/smt/control
